# scripts/digitalocean/attendee-runtime-agent.service
[Unit]
Description=Attendee runtime agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/attendee/runtime-agent.env
ExecStartPre=/usr/bin/test -f /etc/attendee/runtime-agent.env
ExecStart=/usr/local/bin/attendee-runtime-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
